# dev/endpoint/scheduled_task_creation_on_remote_endpoint_using_at.yml
name: Scheduled Task Creation on Remote Endpoint using At
id: 4be54858-432f-11ec-8209-3e22fbd008af
version: 1
date: '2021-11-11'
author: Mauricio Velazco, Splunk
status: production
type: TTP
description: This analytic looks for the execution of `at.exe` with command-line arguments
  utilized to create a Scheduled Task on a remote endpoint. Red Teams and adversaries
  alike may abuse the Task Scheduler for lateral movement and remote code execution.
  The `at.exe` binary internally leverages the AT protocol which was deprecated starting
  with Windows 8 and Windows Server 2012 but may still work on previous versions of
  Windows. Furthermore, attackers may enable this protocol on demand by changing a
  sytem registry key.
data_source:
- Sysmon Event ID 1
search:
  selection1:
    OriginalFileName: at.exe
  selection2:
    Image|endswith: at.exe
  selection3:
    CommandLine: '*\\\\*'
  condition: (selection1 or selection2) and selection3
how_to_implement: To successfully implement this search, you need to be ingesting
  logs with the process name, parent process, and command-line executions from your
  endpoints.
known_false_positives: Administrators may create scheduled tasks on remote systems,
  but this activity is usually limited to a small set of hosts or users.
references:
- https://docs.microsoft.com/en-us/windows-server/administration/windows-commands/at
- https://docs.microsoft.com/en-us/windows/win32/cimwin32prov/win32-scheduledjob?redirectedfrom=MSDN
tags:
  analytic_story:
  - Active Directory Lateral Movement
  - Living Off The Land
  asset_type: Endpoint
  confidence: 60
  impact: 90
  message: A Windows Scheduled Task was created on a remote endpoint from $dest
  mitre_attack_id:
  - T1053
  - T1053.002
  observable:
  - name: dest
    type: Endpoint
    role:
    - Victim
  product:
  - Splunk Enterprise
  - Splunk Enterprise Security
  - Splunk Cloud
  risk_score: 54
  security_domain: endpoint
tests:
- name: True Positive Test
  attack_data:
  - data: https://media.githubusercontent.com/media/splunk/attack_data/master/datasets/attack_techniques/T1053.002/lateral_movement/windows-sysmon.log
    source: XmlWinEventLog:Microsoft-Windows-Sysmon/Operational
    sourcetype: xmlwineventlog
