# ssa_detections/endpoint/ssa___possible_lateral_movement_powershell_spawn.yml
name: Possible Lateral Movement PowerShell Spawn
id: 1151391e-9633-4675-a115-5ec15adcf2ff
version: 1
date: '2024-02-01'
author: Mauricio Velazco, Splunk
type: TTP
status: validation
description: 'The following analytic is designed to identify possible lateral movement
  attacks that involve the spawning of a PowerShell process as a child or grandchild
  process of commonly abused processes. These processes include services.exe, wmiprsve.exe,
  svchost.exe, wsmprovhost.exe, and mmc.exe.\

  Such behavior is indicative of legitimate Windows features such as the Service Control
  Manager, Windows Management Instrumentation, Task Scheduler, Windows Remote Management,
  and the DCOM protocol being abused to start a process on a remote endpoint. This
  behavior is often seen during lateral movement techniques where adversaries or red
  teams abuse these services for lateral movement and remote code execution.'
data_source:
- Windows Security 4688
search: ' $main = from source  | eval timestamp = time  | eval metadata_uid = metadata.uid  |
  eval process_pid = process.pid | eval process_file = process.file | eval process_file_path
  = process_file.path | eval process_file_name = lower(process_file.name) | eval process_cmd_line
  = process.cmd_line | eval actor_user = actor.user | eval actor_user_name = actor_user.name
  | eval actor_process = actor.process | eval actor_process_pid = actor_process.pid
  | eval actor_process_file = actor_process.file | eval actor_process_file_path =
  actor_process_file.path | eval actor_process_file_name = lower(actor_process_file.name)
  | eval device_hostname = device.hostname | where (actor_process_file_name IN ("wmiprvse.exe",
  "services.exe", "svchost.exe", "wsmprovhost.exe", "mmc.exe")) AND ((process_file_name
  IN ("powershell.exe", "pwsh.exe")) OR (process_file_name="cmd.exe" AND (process_cmd_line
  LIKE "%powershell.exe%" OR process_cmd_line LIKE "%pwsh.exe%"))) --finding_report--'
how_to_implement: The detection is based on data that originates from Endpoint Detection
  and Response (EDR) agents. These agents are designed to provide security-related
  telemetry from the endpoints where the agent is installed. To implement this search,
  you must ingest logs that contain the process GUID, process name, and parent process.
  Additionally, you must ingest complete command-line executions. These logs must
  be processed using the appropriate Splunk Technology Add-ons that are specific to
  the EDR product. The logs must also be mapped to the `Processes` node of the `Endpoint`
  data model. Use the Splunk Common Information Model (CIM) to normalize the field
  names and speed up the data modeling process.
known_false_positives: Legitimate applications may spawn PowerShell as a child process
  of the the identified processes. Filter as needed.
references:
- https://attack.mitre.org/techniques/T1021/003/
- https://attack.mitre.org/techniques/T1021/006/
- https://attack.mitre.org/techniques/T1047/
- https://attack.mitre.org/techniques/T1053/005/
- https://attack.mitre.org/techniques/T1543/003/
tags:
  analytic_story:
  - Active Directory Lateral Movement
  - Malicious PowerShell
  - Hermetic Wiper
  - Data Destruction
  - Scheduled Tasks
  asset_type: Endpoint
  confidence: 50
  impact: 90
  mitre_attack_id:
  - T1021
  - T1021.003
  - T1021.006
  - T1047
  - T1053.005
  - T1543.003
  - T1059.001
  - T1218.014
  observable:
  - name: process.pid
    type: Other
  - name: process.file.path
    type: File
  - name: process.file.name
    type: File
  - name: process.cmd_line
    type: Other
  - name: actor.user.name
    type: User Name
  - name: actor.process.pid
    type: Other
  - name: actor.process.file.path
    type: File Name
  - name: actor.process.file.name
    type: File Name
  - name: device.hostname
    type: Hostname
  message: A PowerShell process was spawned as a child process of typically abused
    processes on $dest_device_id$
  product:
  - Splunk Behavioral Analytics
  required_fields:
  - process.pid
  - process.file.path
  - process.file.name
  - process.cmd_line
  - actor.user.name
  - actor.process.pid
  - actor.process.file.path
  - actor.process.file.name
  - device.hostname
  risk_score: 45
  security_domain: endpoint
  mappings:
  - ocsf: process.pid
    cim: process_id
  - ocsf: process.file.path
    cim: process_path
  - ocsf: process.file.name
    cim: process_name
  - ocsf: process.cmd_line
    cim: process
  - ocsf: actor.user.name
    cim: user
  - ocsf: actor.process.pid
    cim: parent_process_id
  - ocsf: actor.process.file.path
    cim: parent_process_path
  - ocsf: actor.process.file.name
    cim: parent_process_name
  - ocsf: device.hostname
    cim: dest
tests:
- name: True Positive Test
  attack_data:
  - data: https://media.githubusercontent.com/media/splunk/attack_data/master/datasets/attack_techniques/T1021.003/lateral_movement/windows_security_xml.log
    source: XmlWinEventLog:Security
    sourcetype: XmlWinEventLog
