# dev_ssa/endpoint/ssa___possible_lateral_movement_powershell_spawn.yml
name: Possible Lateral Movement PowerShell Spawn
id: 22282a2d-dc19-4b88-ac61-6c86ff92904f
version: 1
date: '2024-02-01'
author: Mauricio Velazco, Splunk
status: production
type: TTP
description: 'The following analytic is designed to identify possible lateral movement
  attacks that involve the spawning of a PowerShell process as a child or grandchild
  process of commonly abused processes. These processes include services.exe, wmiprsve.exe,
  svchost.exe, wsmprovhost.exe, and mmc.exe.\

  Such behavior is indicative of legitimate Windows features such as the Service Control
  Manager, Windows Management Instrumentation, Task Scheduler, Windows Remote Management,
  and the DCOM protocol being abused to start a process on a remote endpoint. This
  behavior is often seen during lateral movement techniques where adversaries or red
  teams abuse these services for lateral movement and remote code execution.'
data_source:
- Windows Security 4688
search: 
  selection1:
    actor.process.file.name: 
    - wmiprvse.exe
    - services.exe
    - svchost.exe
    - wsmprovhost.exe
    - mmc.exe
  selection2:
    process.file.name: 
    - powershell.exe
    - pwsh.exe
  selection3:
    process.file.name: cmd.exe
    process.cmd_line|contains: 
    - powershell.exe
    - pwsh.exe
  condition: selection1 and (selection2 or selection3)
how_to_implement: The detection is based on data that originates from Endpoint Detection
  and Response (EDR) agents. These agents are designed to provide security-related
  telemetry from the endpoints where the agent is installed. To implement this search,
  you must ingest logs that contain the process GUID, process name, and parent process.
  Additionally, you must ingest complete command-line executions. These logs must
  be processed using the appropriate Splunk Technology Add-ons that are specific to
  the EDR product. The logs must also be mapped to the `Processes` node of the `Endpoint`
  data model. Use the Splunk Common Information Model (CIM) to normalize the field
  names and speed up the data modeling process.
known_false_positives: Legitimate applications may spawn PowerShell as a child process
  of the the identified processes. Filter as needed.
references:
- https://attack.mitre.org/techniques/T1021/003/
- https://attack.mitre.org/techniques/T1021/006/
- https://attack.mitre.org/techniques/T1047/
- https://attack.mitre.org/techniques/T1053/005/
- https://attack.mitre.org/techniques/T1543/003/
tags:
  analytic_story:
  - Active Directory Lateral Movement
  - Malicious PowerShell
  - Hermetic Wiper
  - Data Destruction
  - Scheduled Tasks
  asset_type: Endpoint
  confidence: 50
  impact: 90
  message: A PowerShell process was spawned as a child process of typically abused
    processes on $dest_device_id$
  mitre_attack_id:
  - T1021
  - T1021.003
  - T1021.006
  - T1047
  - T1053.005
  - T1543.003
  - T1059.001
  - T1218.014
  observable: []
  product:
  - Splunk Behavioral Analytics
  required_fields: []
  risk_score: 45
  security_domain: endpoint
tests:
- name: True Positive Test
  attack_data:
  - data: https://media.githubusercontent.com/media/splunk/attack_data/master/datasets/attack_techniques/T1021.003/lateral_movement/windows_security_xml.log
    source: XmlWinEventLog:Security
    sourcetype: XmlWinEventLog
