# detections/endpoint/scheduled_task_creation_on_remote_endpoint_using_at.yml
name: Scheduled Task Creation on Remote Endpoint using At
id: 4be54858-432f-11ec-8209-3e22fbd008af
version: 1
date: '2021-11-11'
author: Mauricio Velazco, Splunk
type: TTP
datamodel:
- Endpoint
description: This analytic looks for the execution of `at.exe` with command-line arguments
  utilized to create a Scheduled Task on a remote endpoint. Red Teams and adversaries
  alike may abuse the Task Scheduler for lateral movement and remote code execution.
  The `at.exe` binary internally leverages the AT protocol which was deprecated starting
  with Windows 8 and Windows Server 2012 but may still work on previous versions of
  Windows. Furthermore, attackers may enable this protocol on demand by changing a
  sytem registry key.
search: '| tstats `security_content_summariesonly` count min(_time) as firstTime max(_time)
  as lastTime from datamodel=Endpoint.Processes where (Processes.process_name=at.exe
  OR Processes.original_file_name=at.exe) (Processes.process=*\\\\*) by Processes.dest
  Processes.user Processes.parent_process_name Processes.process_name Processes.process
  Processes.process_id Processes.parent_process_id | `drop_dm_object_name(Processes)`
  | `security_content_ctime(firstTime)`| `security_content_ctime(lastTime)` | `scheduled_task_creation_on_remote_endpoint_using_at_filter`'
how_to_implement: To successfully implement this search, you need to be ingesting
  logs with the process name, parent process, and command-line executions from your
  endpoints.
known_false_positives: Administrators may create scheduled tasks on remote systems,
  but this activity is usually limited to a small set of hosts or users.
references:
- https://docs.microsoft.com/en-us/windows-server/administration/windows-commands/at
- https://docs.microsoft.com/en-us/windows/win32/cimwin32prov/win32-scheduledjob?redirectedfrom=MSDN
tags:
  analytic_story:
  - Active Directory Lateral Movement
  - Living Off The Land
  confidence: 60
  context:
  - Source:Endpoint
  - Stage:Lateral Movement
  dataset:
  - https://media.githubusercontent.com/media/splunk/attack_data/master/datasets/attack_techniques/T1053.002/lateral_movement/windows-sysmon.log
  impact: 90
  kill_chain_phases:
  - Exploitation
  message: A Windows Scheduled Task was created on a remote endpoint from $dest
  mitre_attack_id:
  - T1053
  - T1053.002
  observable:
  - name: dest
    type: Endpoint
    role:
    - Victim
  product:
  - Splunk Enterprise
  - Splunk Enterprise Security
  - Splunk Cloud
  required_fields:
  - _time
  - Processes.dest
  - Processes.user
  - Processes.parent_process_name
  - Processes.parent_process
  - Processes.original_file_name
  - Processes.process_name
  - Processes.process
  - Processes.process_id
  - Processes.parent_process_path
  - Processes.process_path
  - Processes.parent_process_id
  risk_score: 54
  security_domain: endpoint
  supported_tas:
  - Splunk_TA_microsoft_sysmon
  asset_type: Endpoint
