# detections/endpoint/schedule_task_with_rundll32_command_trigger.yml
name: Schedule Task with Rundll32 Command Trigger
id: 75b00fd8-a0ff-11eb-8b31-acde48001122
version: 1
date: '2021-04-19'
author: Teoderick Contreras, Splunk
status: production
type: TTP
description: The following analytic detects the creation of suspicious tasks in Windows, specifically tasks using the rundll32 command. It's implemented using Windows Security EventCode 4698 for A scheduled task was created, and looks for tasks executed either via schtasks.exe or TaskService. This behavior is worth identifying as it is commonly used by malware, such as TrickBot, that leverages rundll32 to execute its downloader.

  If a true positive is found, it suggests an attacker is trying to persist within the environment or potentially deliver additional malicious payloads, leading to data theft, ransomware, or other damaging outcomes.

  To implement this analytic, ensure you are ingesting logs with task schedule information from your endpoints. Be aware of potential false positives - legitimate uses of Task Scheduler in your environment may cause benign activities to be flagged.

  Upon triage, review the scheduled task's source and the command to be executed. Capture and inspect any relevant on-disk artifacts, and look for concurrent processes to identify the attack source. This approach helps analysts detect potential threats earlier and mitigate the risks.
data_source:
- Windows Event Log Security 4698
search: '`wineventlog_security` EventCode=4698 | xmlkv Message | search Command IN
  ("*rundll32*") | stats count min(_time) as firstTime max(_time) as lastTime by dest,
  Task_Name, Command, Author, Enabled, Hidden, Arguments | `security_content_ctime(firstTime)`
  | `security_content_ctime(lastTime)` | `schedule_task_with_rundll32_command_trigger_filter`'
how_to_implement: To successfully implement this search, you need to be ingesting
  logs with the task schedule (Exa. Security Log EventCode 4698) endpoints. Tune and
  filter known instances of Task schedule used in your environment.
known_false_positives: unknown
references:
- https://labs.vipre.com/trickbot-and-its-modules/
- https://whitehat.eu/incident-response-case-study-featuring-ryuk-and-trickbot-part-2/
tags:
  analytic_story:
  - Windows Persistence Techniques
  - Trickbot
  - IcedID
  - Living Off The Land
  - Scheduled Tasks
  asset_type: Endpoint
  confidence: 100
  impact: 70
  message: A schedule task process commandline rundll32 arguments $Arguments$ in host
    $dest$
  mitre_attack_id:
  - T1053
  observable:
  - name: dest
    type: Hostname
    role:
    - Victim
  product:
  - Splunk Enterprise
  - Splunk Enterprise Security
  - Splunk Cloud
  required_fields:
  - _time
  - dest
  - Task_Name
  - Command
  - Author
  - Enabled
  - Hidden
  - Arguments
  risk_score: 70
  security_domain: endpoint
tests:
- name: True Positive Test
  attack_data:
  - data: https://media.githubusercontent.com/media/splunk/attack_data/master/datasets/malware/trickbot/tasksched/windows-security.log
    source: WinEventLog:Security
    sourcetype: WinEventLog
