# detections/endpoint/winevent_windows_task_scheduler_event_action_started.yml
name: WinEvent Windows Task Scheduler Event Action Started
id: b3632472-310b-11ec-9aab-acde48001122
version: 1
date: '2021-10-19'
author: Michael Haag, Splunk
type: Hunting
datamodel:
- Endpoint
description: The following hunting analytic assists with identifying suspicious tasks
  that have been registered and ran in Windows using EventID 200 (action run) and
  201 (action completed). It is recommended to filter based on ActionName by specifying
  specific paths not used in your environment. After some basic tuning, this may be
  effective in capturing evasive ways to register tasks on Windows. Review parallel
  events related to tasks being scheduled. EventID 106 will generate when a new task
  is generated, however, that does not mean it ran. Capture any files on disk and
  analyze.
search: '`wineventlog_task_scheduler` EventCode IN ("200","201") | rename ComputerName
  as dest | stats count min(_time) as firstTime max(_time) as lastTime by Message
  dest EventCode category | `security_content_ctime(firstTime)` | `security_content_ctime(lastTime)`
  | `winevent_windows_task_scheduler_event_action_started_filter`'
how_to_implement: Task Scheduler logs are required to be collected. Enable logging
  with inputs.conf by adding a stanza for [WinEventLog://Microsoft-Windows-TaskScheduler/Operational]
  and renderXml=false. Note, not translating it in XML may require a proper extraction
  of specific items in the Message.
known_false_positives: False positives will be present. Filter based on ActionName
  paths or specify keywords of interest.
references:
- https://github.com/redcanaryco/atomic-red-team/blob/master/atomics/T1053.005/T1053.005.md
- https://thedfirreport.com/2021/10/18/icedid-to-xinglocker-ransomware-in-24-hours/
tags:
  analytic_story:
  - IcedID
  - Windows Persistence Techniques
  confidence: 100
  context:
  - Source:Endpoint
  - Stage:Defense Evasion
  dataset:
  - https://media.githubusercontent.com/media/splunk/attack_data/master/datasets/attack_techniques/T1053.005/windows_taskschedule/windows-taskschedule.log
  impact: 80
  kill_chain_phases:
  - Exploitation
  message: A Scheduled Task was scheduled and ran on $dest$.
  mitre_attack_id:
  - T1053.005
  observable:
  - name: dest
    type: Hostname
    role:
    - Victim
  product:
  - Splunk Enterprise
  - Splunk Enterprise Security
  - Splunk Cloud
  required_fields:
  - _time
  - TaskName
  - ActionName
  - EventID
  - dest
  - ProcessID
  risk_score: 80
  security_domain: endpoint
  asset_type: Endpoint
