# detections/endpoint/winevent_scheduled_task_created_within_public_path.yml
name: WinEvent Scheduled Task Created Within Public Path
id: 5d9c6eee-988c-11eb-8253-acde48001122
version: 1
date: '2021-04-08'
author: Michael Haag, Splunk
type: batch
datamodel: []
description: 'The following query utilizes Windows Security EventCode 4698, `A scheduled
  task was created`, to identify suspicious tasks registered on Windows either via
  schtasks.exe OR TaskService with a command to be executed from a user writeable
  file path.\

  The search will return the first time and last time the task was registered, as
  well as the `Command` to be executed, `Task Name`, `Author`, `Enabled`, and whether
  it is `Hidden` or not.\

  schtasks.exe is natively found in `C:\Windows\system32` and `C:\Windows\syswow64`.\

  The following DLL(s) are loaded when schtasks.exe or TaskService is launched -`taskschd.dll`.
  If found loaded by another process, it is possible a scheduled task is being registered
  within that process context in memory.\

  Upon triage, identify the task scheduled source. Was it schtasks.exe or was it via
  TaskService. Review the job created and the Command to be executed. Capture any
  artifacts on disk and review. Identify any parallel processes within the same timeframe
  to identify source.'
search: '`wineventlog_security` EventCode=4698 | xmlkv Message | search Command IN
  ("*\\users\\public\\*", "*\\programdata\\*", "*\\temp\\*", "*\\Windows\\Tasks\\*",
  "*\\appdata\\*") | stats count min(_time) as firstTime max(_time) as lastTime by
  dest, Task_Name, Command, Author, Enabled, Hidden | `security_content_ctime(firstTime)`
  | `security_content_ctime(lastTime)` | `winevent_scheduled_task_created_within_public_path_filter`'
how_to_implement: To successfully implement this search, you need to be ingesting
  Windows Security Event Logs with 4698 EventCode enabled. The Windows TA is also
  required.
known_false_positives: False positives are possible if legitimate applications are
  allowed to register tasks in public paths. Filter as needed based on paths that
  are used legitimately.
references:
- https://research.checkpoint.com/2021/irans-apt34-returns-with-an-updated-arsenal/
- https://www.ultimatewindowssecurity.com/securitylog/encyclopedia/event.aspx?eventID=4698
- https://redcanary.com/threat-detection-report/techniques/scheduled-task-job/
- https://docs.microsoft.com/en-us/windows/win32/taskschd/time-trigger-example--scripting-?redirectedfrom=MSDN
- https://app.any.run/tasks/e26f1b2e-befa-483b-91d2-e18636e2faf3/
tags:
  analytic_story:
  - Windows Persistence Techniques
  - Ransomware
  - Ryuk Ransomware
  automated_detection_testing: passed
  dataset:
  - https://media.githubusercontent.com/media/splunk/attack_data/master/datasets/attack_techniques/T1053.005/taskschedule/windows-security.log
  kill_chain_phases:
  - Privilege Escalation
  mitre_attack_id:
  - T1053.005
  product:
  - Splunk Enterprise
  - Splunk Enterprise Security
  - Splunk Cloud
  required_fields:
  - _time
  - dest
  - Task_Name
  - Description
  - Command
  security_domain: endpoint
  impact: 70
  confidence: 100
  # (impact * confidence)/100
  risk_score: 70
  context:
  - Source:Endpoint
  - Stage:Execution
  - Stage:Persistence
  - Stage:Privilege Escalation
  message: "A windows scheduled task was created (task name=$Task_Name$) on $dest$ by the following command: $Command$"
  observable:
  - name: dest
    type: Endpoint
    role:
    - Victim
  - name: Command
    type: Command
    role:
    - Target
